# deploy/sft.service
[Unit]
Description=Simple File Transfer Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/sft serve --port 12345
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
